# templates/sanoid.service.j2
[Unit]
Description=Sanoid 

[Service]
ExecStart={{ sanoid_cmd_path }} renew
